# etc/supervisor/conf.d/OpenVDM.conf
[program:OVDM_cruise]
command=/usr/bin/python /usr/local/bin/OVDM_cruise.py
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_cruise_STDOUT.log
user=root
autostart=true
autorestart=true
stopsignal=QUIT

[program:OVDM_cruiseDirectory]
command=/usr/bin/python /usr/local/bin/OVDM_cruiseDirectory.py
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_cruiseDirectory_STDOUT.log
user=root
autostart=true
autorestart=true
stopsignal=QUIT

[program:OVDM_dataDashboard]
command=/usr/bin/python /usr/local/bin/OVDM_dataDashboard.py
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_dataDashboard_STDOUT.log
user=root
autostart=true
autorestart=true
stopsignal=QUIT

[program:OVDM_lowering]
command=/usr/bin/python /usr/local/bin/OVDM_lowering.py
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_lowering_STDOUT.log
user=root
autostart=true
autorestart=true
stopsignal=QUIT

[program:OVDM_loweringDirectory]
command=/usr/bin/python /usr/local/bin/OVDM_loweringDirectory.py
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_loweringDirectory_STDOUT.log
user=root
autostart=true
autorestart=true
stopsignal=QUIT

[program:OVDM_md5Summary]
command=/usr/bin/python /usr/local/bin/OVDM_md5Summary.py
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_md5Summary_STDOUT.log
user=root
autostart=true
autorestart=true
stopsignal=QUIT

[program:OVDM_postCollectionSystemTransfer]
command=/usr/bin/python /usr/local/bin/OVDM_postCollectionSystemTransfer.py
process_name=%(program_name)s_%(process_num)s
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_postCollectionSystemTransfer.log
user=root
autostart=true
autorestart=true
stopsignal=INT

[program:OVDM_postDataDashboard]
command=/usr/bin/python /usr/local/bin/OVDM_postDataDashboard.py
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_postDataDashboard.log
user=root
autostart=true
autorestart=true
stopsignal=INT

[program:OVDM_postSetupNewCruise]
command=/usr/bin/python /usr/local/bin/OVDM_postSetupNewCruise.py
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_postSetupNewCruise.log
user=root
autostart=true
autorestart=true
stopsignal=INT

[program:OVDM_postSetupNewLowering]
command=/usr/bin/python /usr/local/bin/OVDM_postSetupNewLowering.py
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_postSetupNewLowering.log
user=root
autostart=true
autorestart=true
stopsignal=INT

[program:OVDM_rebootReset]
command=/usr/bin/python /usr/local/bin/OVDM_rebootReset.py
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_rebootReset_STDOUT.log
user=root
autostart=true
autorestart=false
stopsignal=QUIT

[program:OVDM_runCollectionSystemTransfer]
command=/usr/bin/python /usr/local/bin/OVDM_runCollectionSystemTransfer.py
process_name=%(program_name)s_%(process_num)s
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_runCollectionSystemTransfer_STDOUT.log
user=root
autostart=true
autorestart=true
stopsignal=INT

[program:OVDM_runCruiseDataTransfer]
command=/usr/bin/python /usr/local/bin/OVDM_runCruiseDataTransfer.py
process_name=%(program_name)s_%(process_num)s
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_runCruiseDataTransfer_STDOUT.log
user=root
autostart=true
autorestart=true
stopsignal=INT

[program:OVDM_runShipToShoreTransfer]
command=/usr/bin/python /usr/local/bin/OVDM_runShipToShoreTransfer.py
process_name=%(program_name)s_%(process_num)s
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_runShipToShoreTransfer_STDOUT.log
user=root
autostart=true
autorestart=true
stopsignal=INT

[program:OVDM_scheduler]
command=/usr/bin/python /usr/local/bin/OVDM_scheduler.py
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_scheduler_STDOUT.log
user=root
autostart=true
autorestart=true
stopsignal=QUIT

[program:OVDM_sizeCacher]
command=/usr/bin/python /usr/local/bin/OVDM_sizeCacher.py
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_sizeCacher_STDOUT.log
user=root
autostart=true
autorestart=true
stopsignal=QUIT

[program:OVDM_stopJob]
command=/usr/bin/python /usr/local/bin/OVDM_stopJob.py
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_stopJob_STDOUT.log
user=root
autostart=true
autorestart=true
stopsignal=QUIT

[program:OVDM_testCollectionSystemTransfer]
command=/usr/bin/python /usr/local/bin/OVDM_testCollectionSystemTransfer.py
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_testCollectionSystemTransfer_STDOUT.log
user=root
autostart=true
autorestart=true
stopsignal=QUIT

[program:OVDM_testCruiseDataTransfer]
command=/usr/bin/python /usr/local/bin/OVDM_testCruiseDataTransfer.py
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_testCruiseDataTransfer_STDOUT.log
user=root
autostart=true
autorestart=true
stopsignal=QUIT

[group:collectionSystemTransfers]
programs=OVDM_runCollectionSystemTransfer

[group:cruiseDataTransfers]
programs=OVDM_runCruiseDataTransfer
