# etc/supervisor/conf.d/OVDM_postCollectionSystemTransfer.conf
[program:OVDM_postCollectionSystemTransfer]
command=/usr/bin/python /usr/local/bin/OVDM_postCollectionSystemTransfer.py
process_name=Worker_%(process_num)s
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/OpenVDM/OVDM_postCollectionSystemTransfer.log
stderr_logfile=/var/log/OpenVDM/OVDM_postCollectionSystemTransfer.log
user=root
autostart=true
autorestart=true
stopsignal=INT
